# splitbill-deploy.service
[Unit]
Description=SplitBill GitHub Webhook Deployment Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/www/wwwroot/SplitBillv2/deploy-server
ExecStart=/usr/bin/node deploy-server.js
Restart=on-failure
Environment=NODE_ENV=production
# Environment variables are loaded from .env file

[Install]
WantedBy=multi-user.target
